# hdapsd.upstart
# hdapsd starting script
# laptops seldom have >2 HDD, so we care only about first two

start on device-added SUBSYSTEM=block DEVNAME=sd[ab]
instance $1
expect fork
nice -5

script
	. /etc/sysconfig/hdapsd

	/usr/sbin/hdapsd -b -l -d $1 ${HDAPS_SENSITIVITY:+-s ${HDAPS_SENSITIVITY}} $HDAPSD_OPTIONS
end script
